# .docker/php/supervisord.d/application.ini
[program:consume_action]
command=/home/developer/app/back/yii strategy-consumer/consume-actions
autostart=true
autorestart=true
stderr_logfile=/var/log/superv/consume_action.err.log
stdout_logfile=/var/log/superv/consume_action.out.log
user=developer
numprocs=1
process_name = %(program_name)s_%(process_num)02d

[program:consume_condition]
command=/home/developer/app/back/yii strategy-consumer/consume-conditions
autostart=true
autorestart=true
stderr_logfile=/var/log/superv/consume_condition.err.log
stdout_logfile=/var/log/superv/consume_condition.out.log
user=developer
numprocs=1
process_name = %(program_name)s_%(process_num)02d

[program:consume_quotation]
command=/home/developer/app/back/yii strategy-consumer/consume-quotations
autostart=true
autorestart=true
stderr_logfile=/var/log/superv/consume_quotation.err.log
stdout_logfile=/var/log/superv/consume_quotation.out.log
user=developer
numprocs=1
process_name = %(program_name)s_%(process_num)02d

[program:quotation_checker]
command=node /home/developer/app/back_node/quotation_checker.js
autostart=true
autorestart=true
stderr_logfile=/var/log/superv/quotation_checker.err.log
stdout_logfile=/var/log/superv/quotation_checker.out.log
user=developer
numprocs=1
process_name = %(program_name)s_%(process_num)02d

[program:consume_succeeded_order]
command=/home/developer/app/back/yii strategy-consumer/consume-succeeded-orders
autostart=true
autorestart=true
stderr_logfile=/var/log/superv/consume_succeeded_order.err.log
stdout_logfile=/var/log/superv/consume_succeeded_order.out.log
user=developer
numprocs=1
process_name = %(program_name)s_%(process_num)02d

[program:consume_trailing]
command=/home/developer/app/back/yii strategy-consumer/consume-trailing
autostart=true
autorestart=true
stderr_logfile=/var/log/superv/consume_trailing.err.log
stdout_logfile=/var/log/superv/consume_trailing.out.log
user=developer
numprocs=1
process_name = %(program_name)s_%(process_num)02d

[program:quotations_bittrex]
command=node /home/developer/app/back_node/quotations_bittrex.js
autostart=true
autorestart=true
stderr_logfile=/var/log/superv/quotations_bittrex.err.log
stdout_logfile=/var/log/superv/quotations_bittrex.out.log
user=developer
numprocs=1
process_name = %(program_name)s_%(process_num)02d

[program:instant_action_bittrex]
command=node /home/developer/app/back_node/instant_action_bittrex.js
autostart=true
autorestart=true
stderr_logfile=/var/log/superv/instant_action_bittrex.err.log
stdout_logfile=/var/log/superv/instant_action_bittrex.out.log
user=developer
numprocs=1
process_name = %(program_name)s_%(process_num)02d

[program:instant_action_binance]
command=node /home/developer/app/back_node/instant_action_binance.js
autostart=true
autorestart=true
stderr_logfile=/var/log/superv/instant_action_binance.err.log
stdout_logfile=/var/log/superv/instant_action_binance.out.log
user=developer
numprocs=1
process_name = %(program_name)s_%(process_num)02d

[program:order_checker]
command=/home/developer/app/back/yii order-checker/check-orders
autostart=true
autorestart=true
stderr_logfile=/var/log/superv/order_checker.err.log
stdout_logfile=/var/log/superv/order_checker.out.log
user=developer
numprocs=1
process_name = %(program_name)s_%(process_num)02d

[program:wss]
command=node /app/wss/wss.js
autostart=true
autorestart=true
stderr_logfile=/var/log/superv/wss.err.log
stdout_logfile=/var/log/superv/wss.out.log
user=developer
numprocs=1
process_name = %(program_name)s_%(process_num)02d

[program:wss_pubsub]
command=node /app/wss/pubsub.js
autostart=true
autorestart=true
stderr_logfile=/var/log/superv/wss_pubsub.err.log
stdout_logfile=/var/log/superv/wss_pubsub.out.log
user=developer
numprocs=1
process_name = %(program_name)s_%(process_num)02d

[program:strategy-data]
command=/app/back/yii strategy-consumer/strategy-data
autostart=true
autorestart=true
stderr_logfile=/var/log/superv/strategy_data.err.log
stdout_logfile=/var/log/superv/strategy_data.out.log
numprocs=1
process_name = %(program_name)s_%(process_num)02d

[program:consume_mailout]
command=/app/back/yii consumer-launcher/mailout
autostart=true
autorestart=true
stderr_logfile=/var/log/superv/consume_mailout.err.log
stdout_logfile=/var/log/superv/consume_mailout.out.log
numprocs=1
process_name = %(program_name)s_%(process_num)02d

[program:consume_profit]
command=/app/back/yii strategy-consumer/consume-profit
autostart=true
autorestart=true
stderr_logfile=/var/log/superv/consume_profit.err.log
stdout_logfile=/var/log/superv/consume_profit.out.log
numprocs=5
process_name = %(program_name)s_%(process_num)02d

[program:consume_profit_strategy]
command=/app/back/yii strategy-consumer/consume-profit-strategy
autostart=true
autorestart=true
stderr_logfile=/var/log/superv/consume_profit_strategy.err.log
stdout_logfile=/var/log/superv/consume_profit_strategy.out.log
numprocs=10
process_name = %(program_name)s_%(process_num)02d

[program:bot]
command=node /app/bot/index.js
autostart=true
autorestart=true
stderr_logfile=/var/log/superv/bot.err.log
stdout_logfile=/var/log/superv/bot.out.log
user=developer
numprocs=1
process_name = %(program_name)s_%(process_num)02d
